# unix/supervisor/book_supervisor.ini
[unix_http_server]
file=/tmp/supervisor.sock;

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock;

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[supervisord]
logfile=/tmp/zerobin.log
logfile_maxbytes=50MB
logfile_backups=2
loglevel=trace
pidfile=/tmp/supervisord.pid
nodaemon=false
minfds=1024
minprocs=200
user=root

[program:uwsgi]
command=/root/anaconda3/envs/py3/bin/uwsgi --ini /var/www/book/uwsgi.ini
directory=/var/www/book
environment=PATH='/root/anaconda3/envs/py3/bin:/home/maxime/anaconda3/bin:/home/maxime/anaconda3/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin:/usr/bin/X11'
user=root
autostart=true
autorestart=true
stopsignal=QUIT
