# unix/supervisor/prod_supervisor.ini
[unix_http_server]
file=/tmp/supervisor.sock;

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock;

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[supervisord]
logfile=/tmp/zerobin.log
logfile_maxbytes=50MB
logfile_backups=2
loglevel=trace
pidfile=/tmp/supervisord.pid
nodaemon=false
minfds=1024
minprocs=200
user=root

[program:uwsgi]
command=/root/anaconda/envs/py3/bin/uwsgi --ini /var/www/asso/uwsgi.ini
directory=/var/www/asso
environment=PATH='/root/anaconda/envs/py3/bin:/home/maxime/anaconda/bin:/home/maxime/anaconda/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin:/usr/bin/X11'
user=root
autostart=true
autorestart=true